# sparse/usr/lib/systemd/system/connman.service.d/99-enable-wifi.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/connmand -n -W nl80211 --nobacktrace --config=/etc/connman/main-native.conf
